# 20-rw-access-for-mouse-input.rules
ACTION=="add", SUBSYSTEM=="input", RUN+="/bin/chmod 666 /dev/input/by-path/pci-0000:00:14.0-usb-0:1:1.0-event-mouse"
